# config/systemd/user/restic-maintenance.service
[Unit]
Description=Restic Retention/Prune (NAS)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-%h/.config/restic-backup/env
ExecStart=/usr/bin/env zsh -lc 'source "$HOME/.zshrc.local" >/dev/null 2>&1 || true; exec "$HOME/bin/restic-backup" forget'
Nice=15
IOSchedulingClass=best-effort
IOSchedulingPriority=7

